# ansible/templates/gitea_user_repos/rhde/dev/rhde_image/custom-rpms/usb-automation/99-usb-autoconfig.rules
ACTION=="add", KERNEL=="sd*[!0-9]", SUBSYSTEM=="block", ENV{ID_TYPE}!="floppy", ENV{ID_BUS}=="usb", TAG+="systemd", ENV{USB_DEVICE_PATH}="%E{DEVNAME}", ENV{SYSTEMD_WANTS}="usb-autoconfig.service", RUN+="/usr/bin/usb_check.sh"
